# root/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/execlineb -P
/usr/local/nginx/sbin/nginx -g "daemon off;"