# root/etc/s6-overlay/s6-rc.d/webui/run
#!/command/execlineb -P
cd /webui
node server.js