# root/etc/s6-overlay/s6-rc.d/init-folders/run
#!/bin/ash
mkdir -p \
    /ssl/default \
    /log