# root/etc/s6-overlay/s6-rc.d/init-folders/run
#!/bin/ash
mkdir -p \
    /config/nginx \
    /config/nginx/sites \
    /config/nginx/auth \
    /config/nginx/proxy \
    /config/www \
    /ssl/default \
    /log